# monitoring/memero-monitor.service
[Unit]
Description=MEMERO Monitoring Dashboard
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/memero/monitoring
Environment="PATH=/root/memero/venv/bin:/root/memero/env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Haupt-Befehl
ExecStart=/root/memero/venv/bin/python3 /root/memero/monitoring/monitor.py

# Automatischer Neustart bei Absturz
Restart=always
RestartSec=10

# Timeout-Settings
TimeoutStartSec=30
TimeoutStopSec=30

# Logs zu systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=memero-monitor

# Prozess-Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
